# scripts/hospitalrobot-serial.service
[Unit]
Description=HospitalRobot - probe serial devices and create stable symlinks
After=systemd-udev-settle.service
Wants=systemd-udev-settle.service
DefaultDependencies=yes

[Service]
Type=oneshot
RemainAfterExit=yes
# Đợi tty xuất hiện rồi probe.
ExecStartPre=/bin/sleep 2
ExecStart=/usr/bin/python3 __REPO_DIR__/scripts/init_serial_symlinks.py --required motor line camera
ExecStop=/bin/sh -c 'rm -f /dev/hospitalrobot/*'
TimeoutStartSec=60

[Install]
WantedBy=multi-user.target
